# deploy/systemd/vs-messanger-mvp.service
[Unit]
Description=vs-messanger MVP docker compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/vs-messanger/app
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose --project-name vs-messanger --env-file deploy/.env -f deploy/compose.mvp.yml up -d --build --remove-orphans
ExecStop=/usr/bin/docker compose --project-name vs-messanger --env-file deploy/.env -f deploy/compose.mvp.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
